# Interdrone-Communication/batman-net.service
[Unit]
Description=B.A.T.M.A.N. Mesh Network Setup
After=network.target NetworkManager.service
Wants=network.target

[Service]
Type=oneshot
# Replace the mrrdt-3 with the correct drone name
Environment="PATH=/root/.local/bin:/home/mrrdt-3/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/usr/local/bin/batman-mesh-setup.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

# Set working directory if your script needs to find the uv project
# Replace mrrdt-3 with the correct drone name
WorkingDirectory=/home/mrrdt-3/IARC-10/Interdrone-Communication/

# Run as root (required for network configuration)
User=root

# Restart on failure
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target